# contrib/init/sharecoind.conf
description "Sharecoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env SHARECOIND_BIN="/usr/bin/sharecoind"
env SHARECOIND_USER="sharecoin"
env SHARECOIND_GROUP="sharecoin"
env SHARECOIND_PIDDIR="/var/run/sharecoind"
# upstart can't handle variables constructed with other variables
env SHARECOIND_PIDFILE="/var/run/sharecoind/sharecoind.pid"
env SHARECOIND_CONFIGFILE="/etc/sharecoin/sharecoin.conf"
env SHARECOIND_DATADIR="/var/lib/sharecoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # sharecoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$SHARECOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run sharecoind."
        echo "The setting must appear in $SHARECOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Sharecoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$SHARECOIND_PIDDIR"
    chmod 0755 "$SHARECOIND_PIDDIR"
    chown $SHARECOIND_USER:$SHARECOIND_GROUP "$SHARECOIND_PIDDIR"
    chown $SHARECOIND_USER:$SHARECOIND_GROUP "$SHARECOIND_CONFIGFILE"
    chmod 0660 "$SHARECOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$SHARECOIND_PIDFILE" \
    --chuid $SHARECOIND_USER:$SHARECOIND_GROUP \
    --exec "$SHARECOIND_BIN" \
    -- \
    -pid="$SHARECOIND_PIDFILE" \
    -conf="$SHARECOIND_CONFIGFILE" \
    -datadir="$SHARECOIND_DATADIR" \
    -daemon

